# jobs/weekly.timer
[Unit]
Description=Weekly job timer
Conflicts=6h.service

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
