# jobs/weekly.service
[Unit]
Description=Weekly job service
Conflicts=6h.service

[Service]
Type=oneshot
WorkingDirectory=/home/srsbot/srs
ExecStart=/usr/bin/make history
